# centos6.etc/etc/xdg/autostart/linc-cleanup-sockets.desktop
[Desktop Entry]
Type=Application
Name=ORBit Socket Housekeeping
Exec=/usr/bin/linc-cleanup-sockets
NoDisplay=true
StartupNotify=false
