# centos6.etc/etc/cron.monthly/readahead-monthly.cron
#!/bin/sh
. /etc/sysconfig/readahead
if [ "$READAHEAD_COLLECT" = "yes" ] ; then
	touch /.readahead_collect
fi

